# packages/sx05re/retroarch-sx05re/system.d/tmp-cores.mount
[Unit]
Description=Cores directory
After=storage.mount
After=systemd-tmpfiles-setup.service

[Mount]
What=none
Where=/tmp/cores
Type=overlayfs
Options=lowerdir=/usr/lib/libretro,upperdir=/storage/cores

[Install]
WantedBy=kodi.target
