# packages/sx05re/retroarch-sx05re/system.d/tmp-assets.mount
[Unit]
Description=Assets directory
After=storage.mount
After=systemd-tmpfiles-setup.service

[Mount]
What=none
Where=/tmp/assets
Type=overlayfs
Options=lowerdir=/usr/share/retroarch-assets,upperdir=/storage/assets

[Install]
WantedBy=kodi.target
