# packages/sx05re/retroarch-sx05re/system.d/tmp-joypads.mount
[Unit]
Description=Joypad configs directory
After=storage.mount
After=systemd-tmpfiles-setup.service

[Mount]
What=none
Where=/tmp/joypads
Type=overlayfs
Options=lowerdir=/etc/retroarch-joypad-autoconfig,upperdir=/storage/joypads

[Install]
WantedBy=kodi.target
