# packages/sx05re/retroarch-sx05re/system.d/tmp-shaders.mount
[Unit]
Description=Shaders directory RetroArch
After=systemd-tmpfiles-setup.service

[Mount]
What=none
Where=/tmp/shaders
Type=overlayfs
Options=lowerdir=/usr/share/common-shaders,upperdir=/storage/shaders

[Install]
WantedBy=kodi.target
